# pasteboard.cron
#!/bin/sh

[ ! -z "${MAX}" ] && case ${MAX} in
  *[!0-9]*)
    true  # do nothing
    ;;
  *)
    find "/pasteboard/public/storage" -type f \( -iname \*.jpg -o -iname \*.png -o -iname \*.gif -o -iname \*.jpeg -o -iname \*.svg+xml \) -mtime +${MAX} -exec rm {} \;
    ;;
esac
